# apache-php5/root/etc/service/apache2/run
#!/bin/sh
exec apache2ctl -d /etc/apache2 -f /etc/apache2/apache2.conf -e info -DFOREGROUND